# php/supervisor_php5-fpm.conf
[program:php-fpm]
directory=/
command=/usr/sbin/php5-fpm
user=root
autostart=true
autorestart=true
